# support/bukget_upstart.conf
description "BukGet API"
author "steve@chigeek.com"

start on runlevel [234]
stop on runlevel [0156]

exec /opt/bukget/code/server.py

respawn