# server/bukget_upstart_example.conf
description "BukGet Production Site"
author "steve@chigeek.com"

start on runlevel [234]
stop on runlevel [0156]

exec /path/to/bukget/server/app.py

respawn
